# config/systemd/abkd_app-web.target
[Unit]
After=network.target
Wants=abkd_app-frappe-web.service abkd_app-node-socketio.service

[Install]
WantedBy=multi-user.target